# sv/.config/sv/pipewire-pulse/run
#!/bin/sh

sv check pipewire >/dev/null || exit 1

exec pipewire-pulse 2>&1
