# root/sv/etc/sv/runsvdir-behappier/run
#!/bin/sh

export USER="behappier"
export HOME="/home/behappier"

groups="$(id -Gn "$USER" | tr ' ' ':')"
svdir="$HOME/.local/var/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
