# sv/.config/sv/mpDris2/run
#!/bin/sh

sv check mopidy >/dev/null || exit 1

exec mpDris2 --debug 2>&1
