# sv/.config/sv/emacs/run
#!/bin/sh

export COLORTERM=truecolor

exec emacs --fg-daemon 2>&1
