# sv/.config/sv/mopidy/run
#!/bin/sh

exec mopidy \
    --config "$HOME"/.config/mopidy/mopidy.conf \
    2>&1
