# sv/.config/sv/mopidy/run
#!/bin/sh

exec mopidy \
    --config "$HOME"/.config/mopidy/mopidy.conf \
    >/dev/null 2>&1
